# raspberrypi/build_rpi/services/aliroed-compose.service
[Unit]
Description=AliroEd Docker Service
After=docker.service
Requires=docker.service

[Service]
Restart=always
WorkingDirectory=/aliroed
ExecStart=/usr/bin/docker compose up -d

[Install]
WantedBy=multi-user.target
